# supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=tor -f /etc/tor/torrc

[program:api]
command=python3 tor_proxy_api.py
